# server/infra/systemd/ege-migration-reminder.service
[Unit]
Description=EGE History migration retention reminders
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/ege-history-migration-reminder
