# server/infra/systemd/ege-migration-reminder.timer
[Unit]
Description=Check EGE History Firebase mirror and archive deadlines

[Timer]
OnCalendar=*-*-* 09:00:00 Europe/Moscow
Persistent=true

[Install]
WantedBy=timers.target
